# deploy/ariadne-sync.timer
# Ariadne memfiles sync timer — daily at 04:30 (Linux). Pairs with the service.
[Unit]
Description=Daily Ariadne memory maintenance

[Timer]
OnCalendar=*-*-* 04:30:00
Persistent=true

[Install]
WantedBy=timers.target
